# files/system/etc/systemd/user/first-login-setup.service
[Unit]
Description=First login setup for user applications
ConditionPathExists=!/var/home/%u/.config/first-login-done

[Service]
Type=oneshot
ExecStart=/usr/bin/first-login-setup.sh
ExecStartPost=/usr/bin/touch /var/home/%u/.config/first-login-done
RemainAfterExit=yes

[Install]
WantedBy=default.target
