# files/system/etc/systemd/user/first-login-setup.service
[Unit]
Description=First login setup for user applications
# 检查标记文件是否存在，如果存在则跳过执行
ConditionPathExists=!%h/.config/first-login-done

[Service]
Type=oneshot
# 执行初始化脚本
ExecStart=/usr/bin/first-login-setup.sh
# 成功后创建标记文件
ExecStartPost=/usr/bin/touch %h/.config/first-login-done
RemainAfterExit=yes

[Install]
# 图形界面就绪后，此 target 会被激活
WantedBy=graphical-session.target
